# systemd/system/snap.oracle-cloud-agent.oracle-cloud-agent-updater.service
[Unit]
# Auto-generated, DO NOT EDIT
Description=Service for snap application oracle-cloud-agent.oracle-cloud-agent-updater
Requires=snap-oracle\x2dcloud\x2dagent-20.mount
Wants=network.target
After=snap-oracle\x2dcloud\x2dagent-20.mount network.target snapd.apparmor.service
X-Snappy=yes

[Service]
EnvironmentFile=-/etc/environment
ExecStart=/usr/bin/snap run oracle-cloud-agent.oracle-cloud-agent-updater
SyslogIdentifier=oracle-cloud-agent.oracle-cloud-agent-updater
Restart=on-failure
WorkingDirectory=/var/snap/oracle-cloud-agent/20
TimeoutStopSec=30
Type=simple

[Install]
WantedBy=multi-user.target
